# clock-chime.timer
[Unit]
Description=Clock chimes

[Timer]
OnCalendar=*-*-* *:0/15:*
AccuracySec=1s
Persistent=false

[Install]
WantedBy=default.target
