# clock-chime.timer
[Unit]
Description=Clock chimes

[Timer]
OnCalendar=*-*-* *:0/15:*

[Install]
WantedBy=default.target
