# clock-chime.service
[Unit]
Description=Clock Chimes

[Service]
ExecStart=/home/tomprince/bin/clock-chime
Type=oneshot
