# clock-chime.service
[Unit]
Description=Clock Chimes

[Service]
ExecStart=/usr/bin/clock-chime
Type=oneshot
